# config/vexbot.service
[Unit]
Description=Helper robot

[Service]
Type=dbus
ExecStart=/path/to/vexbot/script 
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
